# deploy/launchd/com.gru.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  com.gru.server — LaunchAgent that supervises the local gru server.

  Install via scripts/install-gru.sh, which renders __HOME__ and __INSTALL__
  before bootstrapping into launchd.

  The deployed binary + frontend live in __INSTALL__ (default
  ~/.local/share/gru), separate from the dev repo. Git hooks under
  scripts/git-hooks/ rebuild and `kickstart -k` this service whenever
  main moves on the deploy host.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.gru.server</string>

  <key>ProgramArguments</key>
  <array>
    <string>__INSTALL__/gru</string>
    <string>server</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__INSTALL__</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>__HOME__</string>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <!-- launchd starts processes with no TERM. gru spawns tmux for sessions,
         which fails with "open terminal failed: terminal does not support
         clear" when TERM is unset. xterm-256color has a complete terminfo
         entry on every modern macOS. -->
    <key>TERM</key>
    <string>xterm-256color</string>
    <!-- launchd also strips locale env. Without LC_CTYPE/LANG, tmux clients
         detect a C locale and downgrade multi-byte glyphs (box-drawing,
         emoji, powerline, nerd-font) to single-byte placeholders before
         they reach the PTY. The -u flag we now pass to tmux is the
         primary defense; LANG here is belt-and-suspenders so any other
         locale-aware tool gru spawns sees UTF-8 too. -->
    <key>LANG</key>
    <string>en_US.UTF-8</string>
    <!-- Pin web/dist to the install dir; otherwise FindWebDist falls back to
         <cwd>/web/dist or <exe-dir>/web/dist, both of which are wrong when
         we run outside the repo. -->
    <key>GRU_WEB_DIST</key>
    <string>__INSTALL__/web/dist</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <!-- Back off on crash loops so a broken build doesn't hammer launchd. -->
  <key>ThrottleInterval</key>
  <integer>10</integer>

  <key>StandardOutPath</key>
  <string>__HOME__/Library/Logs/gru/server.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/Library/Logs/gru/server.log</string>
</dict>
</plist>
